# ops/systemd/sgm-yearly-cleanup.service.template
[Unit]
Description=SGM yearly cleanup (preserves users)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__RUN_USER__
WorkingDirectory=__REPO_ROOT__
Environment=CLEANUP_EXECUTION_MODE=docker
Environment=CLEANUP_CONTAINER_NAME=sgm_api
ExecStart=/bin/bash __REPO_ROOT__/ops/run-yearly-cleanup.sh
Nice=10
IOSchedulingClass=best-effort

[Install]
WantedBy=multi-user.target
