# ops/systemd/sgm-yearly-cleanup.timer
[Unit]
Description=Run SGM yearly cleanup once per year

[Timer]
OnCalendar=*-01-05 03:30:00
Persistent=true
RandomizedDelaySec=15m
Unit=sgm-yearly-cleanup.service

[Install]
WantedBy=timers.target
